# backend/celery.ini
[program:celery-worker]
command=/root/.p38virtualenvs/vueDjango/bin/celery -A backend worker --loglevel=INFO
directory=/www/yue/vue-django-project/vue3-company/backend
user=celeryuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery/worker.log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[program:celery-beat]
command=/root/.p38virtualenvs/vueDjango/bin/celery -A backend beat --loglevel=INFO --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/path/to/project
user=celeryuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery/beat.log